# 2020_10_02-traefik-load-balancing/ansible/roles/fake-service/files/hostname.socket
[Unit]
Description=hostname service

[Socket]
ListenStream=127.0.0.1:8888
Accept=yes

[Install]
WantedBy=sockets.target
